# ha-devops-agent/rootfs/etc/services.d/devops-agent/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# Take down the S6 supervision tree when ha-devops-agent fails
# ==============================================================================

if [[ "$1" -ne 0 ]] && [[ "$1" -ne 256 ]]; then
  bashio::log.warning "DevOps-Agent crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "DevOps-Agent Stopped, restarting"
